# email/garmin-email-bridge.service
[Unit]
Description=Garmin Messenger → Email Bridge (DB watcher + mail sender)
After=network-online.target
Wants=network-online.target

[Service]
User=root
#Group=palchrb
EnvironmentFile=-/etc/default/garmin-bridge-email
ExecStart=/usr/bin/env python3 /usr/local/bin/garmin_bridge_email.py
Restart=always
RestartSec=1

# Nice-ish scheduling
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=6

# Harden a bit
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
